# templates/systemd/stackctl-backup.service
[Unit]
Description=stackctl backup for {{ENV}}
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stackctl backup --env {{ENV}}
