# templates/systemd/stackctl-backup.service
[Unit]
Description=stackctl backup for {{.Env}}
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stackctl backup --env {{.Env}}
