# templates/systemd/stackctl-backup.timer
[Unit]
Description=Nightly stackctl backup timer for {{.Env}}

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=stackctl-backup-{{.Env}}.service

[Install]
WantedBy=timers.target
